# cookbooks/logstash_service/templates/default/logstash.conf.erb
# description "start and stop the logstash"

chdir /usr/share/logstash

exec /usr/bin/java -jar logstash-1.1.13-flatjar.jar agent -f logstash_agent.conf -l logstash.log

start on runlevel [2345]
stop on runlevel [^2345]
